# mediascanner-2.0.conf.in
description "Media Scanner"
author "James Henstridge <james.henstridge@canonical.com>"

start on started dbus

respawn
exec @CMAKE_INSTALL_FULL_LIBDIR@/mediascanner-2.0/mediascanner-service-2.0
